# systemd/jiakuan-tun2socks-gateway.service.example
[Unit]
Description=家宽 SOCKS5 出口网关（tun2socks + socat + AnyTLS + NaiveProxy）
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=<项目根目录>/jiakuan-proxy
ExecStart=<项目根目录>/jiakuan-proxy/scripts/start.sh
ExecStop=<项目根目录>/jiakuan-proxy/scripts/stop.sh
TimeoutStartSec=180
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
